# deck-essentials.service
[Unit]
Description=SteamDeck System Essentials
After=network-online.target
Wants=network-online.target
[Service]
Type=oneshot
User=root
ExecStart=/home/deck/.deck-essentials/deck-essentials.sh
WorkingDirectory=/home/deck/.deck-essentials/
[Install]
WantedBy=multi-user.target
